# vpn-portal.service
# ── Reference file — secrets are injected via: ./deploy.sh init --host x.x.x.x
# ── The actual process is managed by Docker, not systemd
# ── This file is used only to store and transport secrets

[Unit]
Description=rzilient VPN Portal (Docker)
After=network.target wg-quick@wg0.service docker.service

[Service]
Type=simple
User=root
ExecStart=/usr/bin/docker start -a vpn-portal
ExecStop=/usr/bin/docker stop vpn-portal
Restart=on-failure
RestartSec=5

# ── Auto-injected by deploy script ────────────────────────────────────────────
Environment=WG_SERVER_PUBLIC_KEY=<auto_injected>
Environment=WG_SERVER_ENDPOINT=<auto_injected>
Environment=BASE_URL=<auto_injected>
Environment=VPN_SUBNET=10.8.0

# ── Set manually via: ./deploy.sh init --host x.x.x.x ────────────────────────
Environment=GOOGLE_CLIENT_ID=REDACTED.apps.googleusercontent.com
Environment=GOOGLE_CLIENT_SECRET=GOCSPX-REDACTED
Environment=ALLOWED_DOMAINS=rzilient.club,rzilient.tech,rzilient.support
Environment=ADMIN_TOKEN=2f3c46d7deb22aa51aaf6e5ccb46b0a2
Environment=SMTP_HOST=smtp.eu.mailgun.org
Environment=SMTP_PORT=587
Environment=SMTP_USERNAME=postmaster@rzilient.club
Environment=SMTP_PASSWORD=b371fe29c180cad197150200a4ed9094-62916a6c-6377b95d
Environment=SMTP_FROM=VPN Portal tech@rzilient.club

[Install]
WantedBy=multi-user.target