# vpn-portal.service
# ── Reference file — secrets are injected via: ./deploy.sh init --host x.x.x.x
# ── The actual process is managed by Docker, not systemd
# ── This file is used only to store and transport secrets

[Unit]
Description=rzilient VPN Portal (Docker)
After=network.target wg-quick@wg0.service docker.service

[Service]
Type=simple
User=root
ExecStart=/usr/bin/docker start -a vpn-portal
ExecStop=/usr/bin/docker stop vpn-portal
Restart=on-failure
RestartSec=5

# ── Auto-injected by deploy script ────────────────────────────────────────────
Environment=WG_SERVER_PUBLIC_KEY=<auto_injected>
Environment=WG_SERVER_ENDPOINT=<auto_injected>
Environment=BASE_URL=<auto_injected>
Environment=VPN_SUBNET=10.8.0

# ── Set manually via: ./deploy.sh init --host x.x.x.x ────────────────────────
Environment=GOOGLE_CLIENT_ID=<your_google_client_id>
Environment=GOOGLE_CLIENT_SECRET=<your_google_client_secret>
Environment=ALLOWED_DOMAINS=rzilient.club,rzilient.tech
Environment=ADMIN_TOKEN=<your_secure_admin_token>

[Install]
WantedBy=multi-user.target
