# vpn-portal.service
[Unit]
Description=rzilient VPN Portal
After=network.target wg-quick@wg0.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vpn-portal
ExecStart=/opt/vpn-portal/vpn-portal
Restart=on-failure
RestartSec=5

Environment=GOOGLE_CLIENT_ID=<your_google_client_id>
Environment=GOOGLE_CLIENT_SECRET=<your_google_client_secret>
Environment=ALLOWED_DOMAINS=<your_allowed_domains>
Environment=BASE_URL=https://<your_vpn_domain>
Environment=WG_INTERFACE=wg0
Environment=WG_SERVER_PUBLIC_KEY=<auto_injected_by_deploy_script>
Environment=WG_SERVER_ENDPOINT=<auto_injected_by_deploy_script>
Environment=VPN_SUBNET=10.8.0
Environment=STATE_FILE=/etc/wireguard/peers/state.json
Environment=PORT=8080

[Install]
WantedBy=multi-user.target
